# etcd_backup/etcd_backup.service
[Unit]
Description=etcd-backup

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/etcd_backup.sh
