# etcd_backup/etcd_backup.timer
[Unit]
Description=etcd-backup
[Timer]
OnCalendar=*-*-* 5:30:00
Unit=etcd-backup.service
[Install]
WantedBy=timers.target
